# launchd/com.athena.health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Athena OS — daily health-check. $HOME подставляется Слоем 5 bootstrap (sed).
     EnvironmentVariables.PATH задан явно (вкл. $HOME/.local/bin — туда ставится claude) —
     иначе launchd не находит claude/chezmoi (баг status=127). -->
<plist version="1.0">
<dict>
  <key>Label</key><string>com.athena.health</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>$HOME/.claude/scripts/health-check.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>$HOME/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
  <key>StandardOutPath</key><string>$HOME/.claude/health.log</string>
  <key>StandardErrorPath</key><string>$HOME/.claude/health.log</string>
</dict>
</plist>
